# packages/notifier/services/lanflix-monitor.service
[Unit]
Description=LanFlix Daemon
After=network.target

[Service]
# Change and/or create the required user and group.
User=${USER}
Group=${USER}

# The -data=/path argument can be used to force the config/db folder
WorkingDirectory=${SCRIPT_DIR}
ExecStart=/usr/bin/node lanflix-monitor.js

Restart=always
RestartSec=10

StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=lanflix-monitor

[Install]
WantedBy=multi-user.target
